# resources/binaries/var/lib/dpkg/info/cydia.list
/.
/usr
/usr/libexec
/usr/libexec/cydia
/usr/libexec/cydia/firmware.sh
/usr/libexec/cydia/alert
/usr/libexec/cydia/move.sh
/usr/libexec/cydia/exec
/usr/libexec/cydia/asuser
/Applications
/Applications/Cydia.app
/Applications/Cydia.app/manage-up.png
/Applications/Cydia.app/Test
/Applications/Cydia.app/configure.png
/Applications/Cydia.app/Info.plist
/Applications/Cydia.app/folder.png
/Applications/Cydia.app/reload.png
/Applications/Cydia.app/Default.png
/Applications/Cydia.app/Cydia_
/Applications/Cydia.app/unknown.png
/Applications/Cydia.app/changes-up.png
/Applications/Cydia.app/news-dn.png
/Applications/Cydia.app/news-up.png
/Applications/Cydia.app/trusted.png
/Applications/Cydia.app/changes-dn.png
/Applications/Cydia.app/Cydia
/Applications/Cydia.app/pinstripe.png
/Applications/Cydia.app/install-up.png
/Applications/Cydia.app/manage-dn.png
/Applications/Cydia.app/advanced.png
/Applications/Cydia.app/search-dn.png
/Applications/Cydia.app/install-dn.png
/Applications/Cydia.app/search-up.png
/Applications/Cydia.app/icon.png
/var
/var/lib
/var/lib/cydia
/System
/System/Library
/System/Library/LaunchDaemons
/System/Library/LaunchDaemons/com.saurik.Cydia.Firmware
/System/Library/PreferenceBundles
/System/Library/PreferenceBundles/CydiaSettings.bundle
/System/Library/PreferenceBundles/CydiaSettings.bundle/icon.png
